# hack/ansible/led-controller.service.jinja
[Unit]
Description=led-controller

[Service]
ExecStart=/usr/bin/led-controller --led-count {{ led_count }} dashboard -v=debug --homekit-pin {{ homekit_pin }}

[Install]
WantedBy=multi-user.target
